# config/xdg/systemd/user/bedtime.service
[Service]
ExecStart=%E/systemd/user/bedtime
