# config/xdg/systemd/user/bedtime.timer
[Timer]
OnActiveSec=10
OnUnitActiveSec=300

[Install]
WantedBy=timers.target
