# config/xdg/systemd/user/brightness.service
[Unit]
Description=brightness-fluctuator

[Service]
Type=oneshot
ExecStart=%E/systemd/user/brightness
