# config/xdg/systemd/user/wallpaper.path
[Path]
PathModified=%E/bg.png

[Install]
WantedBy=default.target
