# etc/xdg/systemd/user/bedtime.timer
[Timer]
OnActiveSec=10
OnUnitActiveSec=600

[Install]
WantedBy=timers.target
