# etc/xdg/systemd/user/brightness.timer
[Timer]
OnActiveSec=0
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
