# etc/xdg/systemd/user/brightness-fluctuator.service
[Service]
ExecStart=%E/systemd/user/brightness-fluctuator
