# etc/xdg/systemd/user/mail.service
[Service]
ExecStart=%E/systemd/user/mail
