# etc/xdg/systemd/user/mail.timer
[Timer]
OnActiveSec=0
OnUnitActiveSec=300

[Install]
WantedBy=timers.target
