# etc/xdg/systemd/user/mail.service
[Unit]
After==graphical-session.target
PartOf==graphical-session.target

[Service]
ExecStart=%E/systemd/user/mail
