# etc/xdg/systemd/user/uptime.timer
[Timer]
OnActiveSec=0
OnUnitActiveSec=600

[Install]
WantedBy=timers.target
